# deploy/systemd/modelprobe-backend.service
[Unit]
Description=Model Probe Backend
After=network.target docker.service
Wants=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/projects/modelprobe/backend
EnvironmentFile=/opt/projects/modelprobe/backend/.env
ExecStart=/opt/projects/modelprobe/backend/modelprobe-server
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
